# mx-haptics.service.in
[Unit]
Description=MX Master 4 Haptic Feedback Daemon
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=@bindir@/mx4hyprland
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=3
TimeoutStopSec=5

[Install]
WantedBy=graphical-session.target
